# e20-mega-alarm.conf
# Synapse example - start as a service  (copy to /etc/init)
description	"Start e20-kit1-demo web app" 
start on runlevel [2]
stop on runlevel [!2]
exec python /home/snap/mega_alarm/tm_sender.py

